# deploy/systemd/hearthly-stack.service
[Unit]
Description=MakeMenage docker stack
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/hearthly
ExecStart=/usr/bin/env bash /opt/hearthly/scripts/server/install-or-update.sh
ExecStop=/usr/bin/env bash -lc 'cd /opt/hearthly && docker compose --env-file .env.production -f docker-compose.prod.yml down'

[Install]
WantedBy=multi-user.target
